# infrastructure/systemd/metron-intraday.timer
# Every 5 minutes while the trading box is up. The box's lifecycle bounds the
# day; the collector itself double-gates each tick on the NYSE session calendar
# (holidays + half-days via alpha_engine_lib.trading_calendar) AND the Metron UI
# heartbeat — quotes are fetched only while the market is open and the app is
# actually in use, so off-gate ticks cost one S3 GET and exit.
[Unit]
Description=Metron intraday quotes every 5 min (session + app-heartbeat gated in-process)

[Timer]
OnCalendar=*-*-* *:00/5:00
AccuracySec=30
Persistent=false

[Install]
WantedBy=timers.target
