# infrastructure/systemd/metron-intraday.service
# Metron intraday-quotes producer (config#1023) — one shot per timer tick.
# Writes market_data/intraday/latest.json (last / open / prior close per held
# symbol + the major-index ETF proxies, ~15-min delayed) for Metron's Today view
# and Overview markets strip. The collector gates on the US market window; off-session
# ticks exit "skipped" without fetching. The UI-heartbeat demand gate is opt-in —
# append --require-heartbeat below to stop fetching when no Metron app is open
# (multi-tenant cost control; left OFF on this single-tenant owner build so the strip
# stays warm and never freezes at the morning quote after the app has been idle).
# Installed on the trading box by infrastructure/install-metron-intraday.sh;
# S3 write rides the instance role (executor-role already has market_data/*).
[Unit]
Description=Metron intraday quotes producer (alpha-engine-data)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ec2-user
WorkingDirectory=/home/ec2-user/alpha-engine-data
ExecStart=/home/ec2-user/alpha-engine-data/.venv/bin/python -m collectors.metron_market_data --only-intraday --require-heartbeat
TimeoutStartSec=300
